# systemd/blade-id.service
[Unit]
Description=Run the blade-id python script
[Service]
ExecStart=/usr/bin/python /usr/local/bin/blade-id.py
