# systemd/blade-id.service
[Unit]
Description=Run the blade-id python script
[Service]
ExecStart=/usr/bin/blade-id
[Install]
WantedBy=multi-user.target
